# ops/com.pedro.hackathonbtc.status-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pedro.hackathonbtc.status-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-u</string>
    <string>/Users/pedrorios/Library/Application Support/hackathon-btc/status_monitor.py</string>
    <string>--out-dir</string>
    <string>/Users/pedrorios/Library/Application Support/hackathon-btc/status</string>
    <string>--interval-sec</string>
    <string>300</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/pedrorios/Library/Application Support/hackathon-btc</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>STATUS_MONITOR_WORKERS</key>
    <string>8</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/pedrorios/Library/Logs/hackathon-btc/status_monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pedrorios/Library/Logs/hackathon-btc/status_monitor.err.log</string>
</dict>
</plist>
